# iso/systemd_units/luca-shell@.service
# llmbasedos/iso/systemd_units/luca-shell@.service
[Unit]
Description=Luca Shell on %I
Documentation=man:agetty(8) man:systemd-getty-generator(8)
After=systemd-user-sessions.service plymouth-quit-wait.service systemd-logind.service mcp-gateway.service
Before=getty.target
IgnoreOnIsolate=yes
ConditionPathExists=/dev/%I

[Service]
User=llmuser
Group=tty # For TTY access, llmuser must be in tty group
WorkingDirectory=/home/llmuser # Shell starts in user's home
ExecStart=-/usr/bin/python -m llmbasedos.shell.luca
Type=idle
Restart=always # Keep shell running on TTY
RestartSec=1
UtmpIdentifier=%I
TTYPath=/dev/%I
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
KillMode=process-group # Kill shell and its children if service stops
IgnoreSIGPIPE=no
SendSIGHUP=yes
StandardInput=tty
StandardOutput=tty
# StandardError=journal # Errors from shell REPL can go to journal
StandardError=tty # Or directly to TTY for immediate visibility
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONPATH=/opt" # So 'llmbasedos.shell.luca' is found

[Install]
WantedBy=getty.target # This allows 'systemctl enable luca-shell@tty1.service'
