# iso/systemd_units/mcp-agent.service
# llmbasedos/iso/systemd_units/mcp-agent.service
[Unit]
Description=LLMBasedOS MCP Agent Server
After=network-online.target mcp-gateway.service docker.service
Wants=mcp-gateway.service docker.service

[Service]
Type=simple
User=llmuser
Group=llmgroup # Ensure llmuser is also in 'docker' group if using Docker socket directly
ExecStart=/usr/bin/python -m llmbasedos.servers.agent.server
WorkingDirectory=/opt/llmbasedos
Restart=on-failure
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONPATH=/opt"
Environment="LLMBDO_AGENT_WORKFLOWS_DIR=/etc/llmbasedos/workflows"
# Ensure Docker socket is accessible or configure DOCKER_HOST if needed

[Install]
WantedBy=multi-user.target
