# iso/systemd_units/mcp-mail.service
# llmbasedos/iso/systemd_units/mcp-mail.service
[Unit]
Description=LLMBasedOS MCP Mail Server
After=network-online.target mcp-gateway.service
Wants=mcp-gateway.service

[Service]
Type=simple
User=llmuser
Group=llmgroup
ExecStart=/usr/bin/python -m llmbasedos.servers.mail.server
WorkingDirectory=/opt/llmbasedos
Restart=on-failure
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONPATH=/opt"
Environment="LLMBDO_MAIL_ACCOUNTS_CONFIG=/etc/llmbasedos/mail_accounts.yaml"

[Install]
WantedBy=multi-user.target
