# iso/systemd_units/mcp-fs.service
# llmbasedos/iso/systemd_units/mcp-fs.service
[Unit]
Description=LLMBasedOS MCP File System Server
After=network-online.target mcp-gateway.service
Wants=mcp-gateway.service

[Service]
Type=simple
User=llmuser
Group=llmgroup
ExecStart=/usr/bin/python -m llmbasedos.servers.fs.server
WorkingDirectory=/opt/llmbasedos # Python -m resolves from here if PYTHONPATH includes /opt
Restart=on-failure
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONPATH=/opt"
# Environment="LLMBDO_FS_VIRTUAL_ROOT=/home/llmuser/shared_fs" # Example custom root for FS server

[Install]
WantedBy=multi-user.target
